# deploy/lib/systemd/system/sms_poller.service
[Unit]
Description=On boot, poll for any incoming SMS

[Service]
Type=oneshot
ExecStart=/home/new/poll_sms

[Install]
WantedBy=multi-user.target
